# etc/supervisord.d/skydns.ini
[program:skydns]
command=/opt/qnib/skydns/bin/start.sh
autostart=false
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

