# etc/supervisord.d/skydns.ini
[program:skydns]
command=/root/bin/start_skydns.sh
autostart=false
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

